# scripts/jarvis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jeetvaidya.jarvis</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/jeetvaidya/Documents/projects/jarvis/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jeetvaidya/Documents/projects/jarvis</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/jeetvaidya/Documents/projects/jarvis/logs/daemon-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jeetvaidya/Documents/projects/jarvis/logs/daemon-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/jeetvaidya</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
